# Alpine/initd_py-kms
#!/sbin/openrc-run

name="$SVCNAME"
command_background=true
command="/usr/local/share/$SVCNAME/pykms_Server.py"
command_args=":: 1688 -V DEBUG -F /var/log/$SVCNAME.log"
pidfile="/var/run/$SVCNAME.pid"

depend() {
        need net
        after firewall
}
